# daemon/blockpyd.timer
[Unit]
Description=Blockpyd Every hour

[Timer]
Persistent=true
OnCalendar=*:0,15,30,45
OnBootSec=5
Unit=blockpyd.service