# daemon/blockpyd.service

[Unit]
Description=Blockpy service

[Service]
ExecStart=/usr/bin/blockpyd
